# i3/power-limit.service
[Unit]
Description=Power Limit Service

[Service]
ExecStart=/home/hello/.config/i3/power-limit.sh

[Install]
WantedBy=multi-user.target
