# settings_backup/rc.local
#!/bin/sh -e

nvidia-smi -pm 1
nvidia-smi -pl 200
nvidia-smi -gtt 80
nvidia-smi -lgc 210,2100
modprobe v4l2loopback video_nr=0 card_label="v4l2" exclusive_caps=1 max_buffers=2
#modprobe xone-wired

exit 0
